# agent/frameos_agent.service
# /etc/systemd/system/frameos-agent.service
[Unit]
Description=FrameOS Agent (auto-reconnect, hardened)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=%I
WorkingDirectory=/srv/frameos/agent/current
ExecStart=/srv/frameos/agent/current/frameos_agent
Restart=always
RestartSec=5
LimitNOFILE=65536
PrivateTmp=yes
ProtectSystem=full

[Install]
WantedBy=multi-user.target